# fsdb.service
[Unit]
Description=fsdb — filesystem-backed key-value store
After=local-fs.target
Documentation=https://github.com/hamkee/fsdb

[Service]
Type=simple
ExecStart=/usr/local/bin/fsdb
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

User=fsdb
Group=fsdb
LimitNOFILE=65536

RuntimeDirectory=fsdb
StateDirectory=fsdb
LogsDirectory=fsdb

Environment=FSDB_SOCKET=/run/fsdb/fsdb.sock
Environment=FSDB_DBDIR=/var/lib/fsdb
Environment=FSDB_LOGDIR=/var/log/fsdb
Environment=FSDB_PIDFILE=/run/fsdb/fsdb.pid

ProtectSystem=strict
ReadWritePaths=/run/fsdb /var/lib/fsdb /var/log/fsdb
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
NoNewPrivileges=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
